# pirobot.service
[Unit]
Description=PiRobot V4 - Raspberry Pi 5 Optimized
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/home/pi/PiRobot5/PiRobot/bin/python main.py
WorkingDirectory=/home/pi/PiRobot5
User=pi
Group=pi
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/PiRobot5
Environment=PYTHONOPTIMIZE=1
Environment=OMP_NUM_THREADS=4
Environment=MKL_NUM_THREADS=4
Restart=always
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=30
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pirobot
# Pi 5 specific optimizations
Nice=-5
IOSchedulingClass=1
IOSchedulingPriority=4

[Install]
WantedBy=multi-user.target 